# bot-paris.service
[Unit]
Description=Bot Discord - Paris personnalises
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

# ⚠️ À ADAPTER : remplace UTILISATEUR (commande `whoami`) et vérifie le
# chemin du dossier ainsi que celui de node (commande `which node`).
User=UTILISATEUR
WorkingDirectory=/home/UTILISATEUR/GOONSOCIETY
ExecStart=/usr/bin/node index.js

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
